# rootfs/etc/services.d/ems/run
#!/usr/bin/with-contenv bash
set -euo pipefail
cd /app
exec python3 server.py
